# Templates/RhelAppDemo.service.j2
[Unit]
Description=RhelAppDemo service
After=network.target

[Service]
WorkingDirectory=/home/ec2-user/DemoApp
ExecStart=npm start